# com.taomonitor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.taomonitor</string>
  <key>ProgramArguments</key>
  <array>
    <!-- Edit this path to match your actual venv and project location -->
    <string>/Users/YOUR_USERNAME/dev/tao-subnet-investigation/.venv/bin/python</string>
    <string>/Users/YOUR_USERNAME/dev/tao-subnet-investigation/main.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USERNAME/dev/tao-subnet-investigation</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/dev/tao-subnet-investigation/logs/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/dev/tao-subnet-investigation/logs/launchd.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
